# apps/webui/deploy/systemd/trusttunnel-webui.service
[Unit]
Description=TrustTunnel WebUI
After=network.target

[Service]
Type=simple
User=trusttunnel
Group=trusttunnel
WorkingDirectory=/opt/trusttunnel-suite/apps/webui
EnvironmentFile=/opt/trusttunnel-suite/apps/webui/deploy/systemd/trusttunnel-webui.env
ExecStart=/opt/trusttunnel-suite/apps/webui/build/trusttunnel-webui
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
